# systemd/ledmon.service.in
[Unit]
Description=Enclosure LED Utilities

[Install]
WantedBy=multi-user.target

[Service]
Type=forking
User=root
ExecStart=@sbindir@/ledmon
Restart=on-failure
